# pi-scanner/deploy/careflow-kiosk.service
# systemd *user* unit (NOT a system unit). Install to ~/.config/systemd/user/ and launch it from the
# labwc XDG autostart entry — a user unit runs in the same graphical session as the compositor, so
# Chromium attaches to the Wayland socket. (A system unit would land in a different cgroup/PAM/env and
# fail to attach.) The matching autostart entry + labwc window rules live outside this file; see
# deploy/README.md.
[Unit]
Description=Careflow Scanner Kiosk (Chromium borderless app window on labwc/Wayland)
Documentation=file:///home/pi/scanner/kiosk.sh
# The FastAPI backend must answer before Chromium opens (kiosk.sh also waits, this makes the
# dependency explicit for restarts).
After=careflow-scanner.service
Wants=careflow-scanner.service

[Service]
Type=simple
# Clear any stray Chromium from an earlier (unsupervised) autostart launch so we never end up with
# two overlapping kiosk windows. `-` = ignore "no process matched" on a clean restart. Run pkill
# directly (not via sh -c) so its own command line can't match the pattern. Then let it settle so
# the Wayland surface and Chromium singleton socket are released before kiosk.sh relaunches.
ExecStartPre=-/usr/bin/pkill -f class=careflow-kiosk
ExecStartPre=/bin/sleep 1
ExecStart=/home/pi/scanner/kiosk.sh
Restart=always
RestartSec=3

[Install]
WantedBy=graphical-session.target
