# pi-scanner/deploy/careflow-scanner.service
[Unit]
Description=Careflow RS-485 Scanner Station (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/rak3112-rs485-node/pi-scanner
# Firmware-build service (this repo's api/, or the hub after cutover) + DEV ChirpStack only.
Environment=SCANNER_API_BASE=http://10.10.8.169:8000
Environment=CS_BASE=http://10.10.8.140:8080
# Session bus so the kiosk on-screen-keyboard toggle (squeekboard, sm.puri.OSK0) works from here.
# Assumes the desktop user is uid 1000; adjust if not.
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
# SCANNER_API_TOKEN / node port can be added here or in an EnvironmentFile.
ExecStart=/home/pi/rak3112-rs485-node/pi-scanner/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8080
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
